# image_chroot/etc/service/remoteapp/run
#!/bin/sh
cd /tmp
python /root/client.py
